# server/rmtc.service
[Unit]
Description=RMTC
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rmtcd

[Install]
WantedBy=multi-user.target
